# files/deepclaude/headroom.service
[Unit]
Description=Headroom AI Proxy — context compression for Claude
Documentation=https://github.com/nicholasgriffintn/headroom-ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/headroom proxy --memory --learn --code-aware --anthropic-api-url https://api.deepseek.com/anthropic
Restart=on-failure
RestartSec=5
Environment=HEADROOM_HOST=127.0.0.1
Environment=HEADROOM_PORT=8787

# --- Security hardening ---
# Filesystem
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.headroom
PrivateTmp=yes
PrivateDevices=yes
# Privileges
NoNewPrivileges=yes
# Network (only AF_INET + AF_INET6 + AF_UNIX)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
# Resources
MemoryMax=1G
# Sandbox (allow exec so headroom can spawn itself for workers)
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
SystemCallArchitectures=native
RestrictRealtime=yes

[Install]
WantedBy=default.target
