# files/deepclaude/headroom.service
[Unit]
Description=Headroom AI Proxy — context compression for Claude
Documentation=https://github.com/nicholasgriffintn/headroom-ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/headroom proxy --memory --learn --code-aware --anthropic-api-url https://api.deepseek.com/anthropic
Restart=on-failure
RestartSec=5
Environment=HEADROOM_HOST=127.0.0.1
Environment=HEADROOM_PORT=8787

[Install]
WantedBy=default.target
